# .config/systemd/user/dailycommit.service
[Unit]
Description=Daily commit for selected git repos

[Service]
Type=oneshot
IOSchedulingClass=idle
ExecStart=%h/.local/bin/dailycommit