# .config/systemd/user/dailysync.timer
[Unit]
Description=Daily sync for selected git folders

[Timer]
OnBootSec=15min
OnUnitActiveSec=1d

[Install]
WantedBy=default.target
