# .config/systemd/user/dailysync.service
[Unit]
Description=Daily sync for selected folders

[Service]
Type=oneshot
IOSchedulingClass=idle
ExecStart=%h/.local/bin/dailysync